# profiles/prometheus/uci-defaults/90-firewall-config
uci -q batch << EOI

add firewall zone
set firewall.@zone[-1].name='captive'
set firewall.@zone[-1].network='captive'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'

add firewall zone
set firewall.@zone[-1].name='encrypted'
set firewall.@zone[-1].network='encrypted'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'

add firewall zone
set firewall.@zone[-1].name='ports'
set firewall.@zone[-1].network='ports'
set firewall.@zone[-1].input='ACCEPT'
set firewall.@zone[-1].output='ACCEPT'
set firewall.@zone[-1].forward='ACCEPT'

add firewall forwarding
set firewall.@forwarding[-1].src='captive'
set firewall.@forwarding[-1].dest='wan'

add firewall forwarding
set firewall.@forwarding[-1].src='encrypted'
set firewall.@forwarding[-1].dest='wan'

add firewall forwarding
set firewall.@forwarding[-1].src='ports'
set firewall.@forwarding[-1].dest='wan'

commit firewall

EOI