# profiles/artemisa/uci-defaults/99-openwisp-onboard
rm /etc/config/openwisp
touch /etc/config/openwisp

uci batch << EOI
add openwisp controller
rename openwisp.@controller[0]=http
set openwisp.http.url="https://controller.wayru.net"
set openwisp.http.shared_secret="yvHGMsfHVScJcVAYtcm07oG2TKz2l1HH"
set openwisp.http.verify_ssl="1"
set openwisp.http.management_interface="tun0"
set openwisp.http.tags="wayru-services v23 radsec passpoint-ap"
commit openwisp
EOI