# profiles/helios/uci-defaults/90-dhcp-config
uci -q batch << EOI

set dhcp.captive=dhcp
set dhcp.captive.interface='captive'
set dhcp.captive.start='100'
set dhcp.captive.limit='150'
set dhcp.captive.leasetime='12h'
set dhcp.captive.force='1'

set dhcp.encrypted=dhcp
set dhcp.encrypted.interface='encrypted'
set dhcp.encrypted.start='100'
set dhcp.encrypted.limit='150'
set dhcp.encrypted.leasetime='12h'
set dhcp.encrypted.force='1'

set dhcp.ports=dhcp
set dhcp.ports.interface='ports'
set dhcp.ports.start='100'
set dhcp.ports.limit='150'
set dhcp.ports.leasetime='12h'
set dhcp.ports.force='1'

del dhcp.lan

commit dhcp

EOI